# ansible/roles/membership-kiosk/templates/upcoming-events.service.j2
[Unit]
Description=Upcoming Events Server
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/upcoming-events
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
